# systemd/monitoring_glances/glances.service
[Unit]
Description = Glances web server
After = network.target
 
[Service]
ExecStart = /usr/bin/glances -w --theme-black --process-short-name --hide-kernel-threads --fs-free-space
 
[Install]
WantedBy=multi-user.target
